# crontab
0 3 * * * /root/.acme.sh/acme.sh --renew --dns dns_aws >> /var/log/cron.log 2>&1
0 */6 * * * /root/.acme.sh/deploy/deploy_all_from_acme.sh >> /var/log/cron.log 2>&1

